# surfcode/hotspot_premium/hotspot_premium.service
[Unit]
Description=Hotspot Daemon

[Service]
ExecStart=/apps/hotspot_premium/hotspot_premium.sh start

[Install]
WantedBy=multi-user.target
