# hotspot_daemon/hotspot.service
[Unit]
Description=Hotspot Daemon

[Service]
ExecStart=/apps/hotspot/hotspot.sh start

[Install]
WantedBy=multi-user.target
